# deploy/systemd/lynxiglam-web.service
# /etc/systemd/system/lynxiglam-web.service
[Unit]
Description=LynxiGlam storefront (Next.js 16 standalone / Node 24)
Documentation=file:///opt/lynxiglam/RUNBOOK.md
After=network-online.target lynxiglam-api.service nginx.service
Wants=network-online.target

# StartLimit* belong to [Unit], not [Service] (systemd moved them in v229).
# Left in [Service] they are parsed as an unknown key and silently dropped,
# so the defaults apply: 5 starts per 10s. With RestartSec=5s a crash loop
# only manages ~2 starts per 10s window, so the limit can NEVER trip and a
# broken unit flaps forever instead of giving up — filling the shared
# journald ring and hammering MySQL on every restart.
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=exec
Slice=lynxiglam.slice
User=lynxiglam-web
Group=lynxiglam-web
WorkingDirectory=/opt/lynxiglam/current/web

Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
# V8 must GC before the cgroup kills us: 384m old-space + ~100m off-heap < 512m.
Environment=NODE_OPTIONS=--max-old-space-size=384
Environment=TZ=UTC
# NOTE: NEXT_PUBLIC_* is NOT set here on purpose — it is inlined at BUILD time.
# Setting it at runtime does nothing; see the RUNBOOK build step.

ExecStart=/usr/bin/node /opt/lynxiglam/current/web/server.js

Restart=always
RestartSec=5s
TimeoutStopSec=20s
KillSignal=SIGTERM

MemoryAccounting=yes
MemoryHigh=448M
MemoryMax=512M
MemorySwapMax=0
CPUAccounting=yes
CPUQuota=150%
CPUWeight=50
TasksMax=128
LimitNOFILE=4096

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectProc=invisible
ProtectClock=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RemoveIPC=yes
UMask=0077
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
MemoryDenyWriteExecute=no

# Next needs .next writable, not just .next/cache. Only the FETCH cache lives
# under .next/cache; app-router ISR pages go to .next/server/app —
# file-system-cache.js maps APP_PAGE/APP_ROUTE to path.join(serverDistDir,
# 'app', pathname), and serverDistDir is .next/server.
#
# With only .next/cache whitelisted, ProtectSystem=strict kept .next/server
# read-only, so all 21 `revalidate = 300` routes (/pages/* x15, /blogs/news +
# 5 posts, /pages/reviews) re-rendered every 5 minutes and then failed to
# persist with EROFS. Next swallows that and only console.warn()s, so nothing
# 500s — it just writes one warning per route per 5 minutes into the SHARED
# journald ring forever (pushing out the neighbours' history), serves from the
# in-process memory cache (up to 50MB of the 512MB budget), and snaps back to
# build-time HTML on every restart.
ReadWritePaths=/opt/lynxiglam/current/web/.next

# SSR resolves lynxiglam.kejing.online -> 127.0.0.1 via /etc/hosts (RUNBOOK
# step 14), so loopback is all it ever needs. If you skip the /etc/hosts entry,
# you must add IPAddressAllow=<host public IP>/32 here.
IPAddressDeny=any
IPAddressAllow=localhost

[Install]
WantedBy=multi-user.target
